# etc/systemd/system/starlight-update-issue.service
[Unit]
Description=Update /etc/issue with host IP for Starlight
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/scripts/update-issue.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
